# test/etc/init.d/pygrim-test
### BEGIN INIT INFO
# Provides:          pygrim-test
# Required-Start:    
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: pygrim-test
# Description:       This file should be used to construct scripts to be
#                    placed in /etc/init.d.
### END INIT INFO
# /etc/init.d/uwsgi
#

daemon=/usr/local/bin/uwsgi
pid=/www/pygrim-test/run/pygrim-test.pid
CONF="/www/pygrim-test/conf/pygrim-test.yaml"
args="--yaml $CONF"
pname="pygrim-test_m"

# Carry out specific functions when asked to by the system
case "$1" in
    start)
        echo "Starting uwsgi"
        start-stop-daemon -p $pid --start --exec $daemon -- $args
        ;;
    stop)
        echo "Stopping script uwsgi"
        pidof $pname > $pid
        start-stop-daemon --signal INT -p $pid --stop $daemon --retry 10 -- $args
        ;;
    restart)
        echo "Stopping script uwsgi"
        pidof $pname > $pid
        start-stop-daemon --signal INT -p $pid --stop $daemon --retry 10 -- $args
        sleep 1
        echo "Starting uwsgi"
        start-stop-daemon -p $pid --start --exec $daemon -- $args
        ;;
    reload)
        echo "Reloading conf"
        pidof $pname > $pid
        kill -HUP $(cat $pid)
        ;;
    *)
        echo "Usage: /etc/init.d/uwsgi {start|stop|restart|reload}"
        exit 1
    ;;
esac

exit 0

